# Server/Root/etc/systemd/system/BackupAll.service
[Unit]
Description=BackupAll.service
After=network.target

[Service]
User=root
ExecStart=/usr/bin/bash -c "zx /Main/Server/Scripts/BackupAll.zx.mjs"
