# Server/Root/etc/systemd/system/FreshRSS-actualize.service
[Unit]
Description=FreshRSS get new content
Wants=freshrss.timer

[Service]
User=www-data
Type=simple
ExecStart=/usr/bin/php /Server/www/FreshRSS/app/actualize_script.php
