# Server/Root/etc/systemd/system/FreshRSS-actualize.timer
[Unit]
Description=FreshRSS get new content

[Timer]
OnBootSec=30s
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
